# packaging/vegad/vegad-update-check.service
[Unit]
Description=Vega periodic update check

[Service]
Type=oneshot
ExecStart=/usr/lib/vega/vegad check-updates
