# packaging/debian-src/debian/vegad.prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ] && [ -d /run/systemd/system ]; then
	systemctl disable --now vegad-update-check.timer 2>/dev/null || true
fi

#DEBHELPER#

exit 0
